# tools/network/ghost-prefer-phone-hotspot.service
[Unit]
Description=Prefer GhostChessboard phone hotspot when visible
Wants=NetworkManager.service
After=NetworkManager.service

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/ghost-prefer-phone-hotspot.sh
